# etc/profile.d/aliases.sh
# Setup some handy aliases
alias kube-system='kubectl --namespace=kube-system'
alias default='kubectl --namespace=default'


