# rootfs/etc/profile.d/banner.sh
#!/usr/bin/env bash
COLOR_RESET="[0m"
BANNER_INDENT="${BANNER_INDENT:-  }"
BANNER_COLOR="${BANNER_COLOR:- [36m}"