# rootfs/etc/profile.d/direnv.sh
# Install `direnv` via PROMPT_COMMAND hook
#

eval "$(direnv hook bash)"
