# rootfs/etc/profile.d/geodesic.sh
export GEODESIC_SHELL=true
export SSH_AGENT_CONFIG=${LOCAL_STATE}/.ssh-agent


