# rootfs/etc/profile.d/geodesic.sh
export GEODESIC_SHELL=true
export SSH_AGENT_CONFIG=/var/tmp/.ssh-agent
