# rootfs/etc/profile.d/geodesic.sh
#!/usr/bin/env bash
export GEODESIC_SHELL=true
export TF_VAR_aws_assume_role_arn=${AWS_ROLE_ARN}