# rootfs/etc/profile.d/syslog-ng.sh
#!/usr/bin/env bash
syslog-ng -f /etc/syslog-ng/syslog-ng.conf
