# rootfs/etc/profile.d/syslog-ng.sh
syslog-ng -f /etc/syslog-ng/syslog-ng.conf
