# rootfs/etc/profile.d/tsh-login.sh
#!/bin/bash

function tsh-login() {
	tsh login --bind-addr=:${GEODESIC_PORT} --proxy=${TELEPORT_PROXY_DOMAIN_NAME:-tele.${DOCKER_IMAGE#*.}} $*
}
