# rootfs/etc/profile.d/use-profile.sh
if [ -n "${AWS_PROFILE}" ]; then
  use-profile
fi
